# trading-system.service
[Unit]
Description=Autonomous Trading System Daemons
After=network.target docker.service
Requires=docker.service

[Service]
Type=forking
User=linuxuser
WorkingDirectory=/home/linuxuser/trading-system
EnvironmentFile=/home/linuxuser/.trading_env
ExecStart=/home/linuxuser/trading-system/start_trading_system.sh
ExecStop=/home/linuxuser/trading-system/start_trading_system.sh --stop
RemainAfterExit=yes
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
